# home_visit_profit_bot/deploy/homevisit-api.service
# systemd-юнит для backend REST API — на uvicorn/FastAPI (Фаза 1).
#
# Установка / обновление:
#   sudo cp deploy/homevisit-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl restart homevisit-api
#
# ОТКАТ на старый ThreadingHTTPServer: раскомментировать ExecStart с app.main,
# закомментировать uvicorn, daemon-reload + restart. Старый entrypoint (app/main.py)
# оставлен специально ради отката — он поднимает тот же REST на том же порту.

[Unit]
Description=Home Visit backend REST API (FastAPI/uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=homevisit
Group=homevisit
WorkingDirectory=/opt/tba/home_visit_profit_bot
EnvironmentFile=/opt/tba/home_visit_profit_bot/.env

# Лимиты (Фаза 1.7): воркеров немного (сервер 1 маленький), всплеск гасит
# --limit-concurrency (сверх лимита — 503, а не рост памяти без предела); короткий
# keep-alive не даёт простаивающим соединениям копиться. Размер тела режет Caddy.
Environment=DB_POOL_MIN=2
Environment=DB_POOL_MAX=20

# Миграция схемы ДО старта воркеров. Старый сервер (app.main) звал init_db при
# старте, uvicorn-путь этого не делал — 17.07.2026 это уронило приём GPS на бою:
# код ссылался на новую колонку left_at, а ALTER никто не выполнил (UndefinedColumn,
# 500 на каждой точке вне геозоны). init_db идемпотентен; отдельный шаг исключает
# и гонку двух воркеров за ALTER, если бы миграция жила в lifespan.
ExecStartPre=/opt/tba/home_visit_profit_bot/.venv/bin/python -c 'from app.config import load_config; from app.db import init_db; init_db(load_config())'

ExecStart=/opt/tba/home_visit_profit_bot/.venv/bin/uvicorn app.asgi:app --host 127.0.0.1 --port 8088 --workers 2 --limit-concurrency 64 --timeout-keep-alive 15 --no-access-log

# ОТКАТ (закомментировано):
# ExecStart=/opt/tba/home_visit_profit_bot/.venv/bin/python -m app.main

Restart=on-failure
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tba/home_visit_profit_bot

[Install]
WantedBy=multi-user.target
